# deploy/systemd/crypto-trader.service
[Unit]
Description=Crypto Trader multi-symbol trading daemon
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=notify
NotifyAccess=main
Environment=PROJECT_DIR=%h/workspace/crypto-trader
Environment=TZ=Asia/Seoul
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=%h/workspace/crypto-trader
ExecStartPre=/usr/bin/bash -lc 'mkdir -p "$PROJECT_DIR/artifacts"'
ExecStart=/usr/bin/bash -lc 'cd "$PROJECT_DIR" && exec "$PROJECT_DIR/.venv/bin/crypto-trader" run-multi --config config/daemon.toml >> "$PROJECT_DIR/artifacts/daemon.log" 2>&1'
Restart=always
RestartSec=15
WatchdogSec=180
TimeoutStopSec=45
KillSignal=SIGTERM
SyslogIdentifier=crypto-trader

[Install]
WantedBy=default.target
